# deploy/ansible/roles/space-inference-gateway/templates/com.slushsystems.space-inference-gateway.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>{{ gateway_label }}</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{{ gateway_src_dir }}/run-gateway.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>{{ ansible_facts.env.HOME }}/Library/Logs/space-inference-gateway.log</string>
  <key>StandardErrorPath</key><string>{{ ansible_facts.env.HOME }}/Library/Logs/space-inference-gateway.log</string>
</dict>
</plist>
